# Documents/Linux/super-duper-octo-spork.service
[Unit]
After=network.target

[Service]
ExecStart=sudo /Documents/Linux/Update.sh

[Install]
WantedBy=multi-user.target
